# deployment/car-logger-restart.timer
[Unit]
Description=Restart car-logger every day at 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
